# crontab
* * * * * root /app/script.sh
# 需要在这个文件的最后留一个空白行
# 貌似cron最小也就只能每分钟执行一次，可以在script.sh里面写个循环做到按秒执行
# 或者还有个项目  https://github.com/mcuadros/ofelia, 更加灵活

